# backend/Scripts/crontab
0 3 * * * python /home/cleanupCrew.py >> /logfile